# gtr5/systemd/memo-update.timer
[Unit]
Description=Daily trigger memo-update

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target


