# gtr5/systemd/memo-update.service
[Unit]
Description=memo git sync update
After=org.gnome.Shell.target

[Service]
ExecStart=/home/fool/bin/memo-update


