# pi/systemd/daily-fetch.service
[Unit]
Description=Daily fetch keep up-to-date repositories (worker)
After=network.target

[Service]
Type=idle
User=pi
ExecStart=/home/pi/.local/bin/daily-fetch.sh
