# pi/systemd/daily-fetch.timer
[Unit]
Description=Daily fetch keep up-to-date repositories

[Timer]
# TODO

[Install]
WantedBy=multi-user.target
