# pi/systemd/daily-fetch.timer
[Unit]
Description=Daily fetch keep up-to-date repositories

[Timer]
OnCalendar=*-*-* 13:00:00

[Install]
WantedBy=multi-user.target
