# pi/systemd/portr-send-ka.timer
[Unit]
Description=Every minute send alive message to web server

[Timer]
OnStartupSec=30
OnUnitActiveSec=60

[Install]
WantedBy=multi-user.target
